# assets/kde/config/autostart/editor-demo.desktop
[Desktop Entry]
Type=Application
Exec=/home/user/bin/editor-demo
